# puppet/modules/services_website/files/services-website.service
[Unit]
Description=Coding-Tag website 

[Service]
Type=simple
ExecStart=/usr/local/bin/flask run
Environment=FLASK_APP=/path/website.py

[Install]
WantedBy=multi-user.target
